# supervisord.conf
[supervisord]
nodaemon=true

[program:golang]
command=/app/binary
